# my_program_supervisor.ini
[program:my_program]
command=python3 /Users/gedizucar/Desktop/main.py
directory=/Users/gedizucar/Desktop/
autostart=true
autorestart=true
stderr_logfile=/var/log/my_program.err.log
stdout_logfile=/var/log/my_program.out.log

